# dot_config/systemd/user/night-light-off.service
[Unit]
Description=Turn off Night Light (Reset Gamma)

[Service]
Type=oneshot
ExecStart=pkill hyprsunset
