# dot_config/systemd/user/night-light-off.timer
[Unit]
Description=Timer to turn off Night Light

[Timer]
OnCalendar=07:00
Persistent=true

[Install]
WantedBy=timers.target
