# dot_config/systemd/user/night-light-on.service
[Unit]
Description=Turn on Night Light (Gammastep)

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprsunset -t 4000
