# dot_config/systemd/user/night-light-on.timer
[Unit]
Description=Timer to turn on Night Light

[Timer]
OnCalendar=22:00
Persistent=true

[Install]
WantedBy=timers.target
