# dot_config/systemd/user/transformers_ocr.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/transformers_ocr start --foreground
