# files/.config/systemd/user/hyprsunset.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/bash -c 'exec /usr/bin/hyprsunset 2>&1 | stdbuf -oL grep -v "\\[TRACE\\]"'
KillMode=control-group
